# udev/99-zmirror.rules
SUBSYSTEM=="block", RUN+="/usr/local/bin/zmirror-trigger"